# gpu-services/jarvis-whisper-stt.service
[Unit]
Description=Jarvis Voice Whisper STT Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/jarvis-voice
EnvironmentFile=%h/jarvis-voice/.env
ExecStart=%h/jarvis-voice/venv/bin/python3 %h/jarvis-voice/gpu-services/whisper_stt_service.py --model ${WHISPER_MODEL}
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
