# gpu-services/jarvis-piper-tts.service
[Unit]
Description=Jarvis Piper TTS Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/jarvis-voice
ExecStart=%h/jarvis-voice/venv/bin/python3 %h/jarvis-voice/gpu-services/piper_tts_service.py
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
