# gpu-services/jarvis-speaker-verify.service
[Unit]
Description=Jarvis Speaker Verification Service
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=%h/jarvis-voice
ExecStart=%h/jarvis-voice/venv/bin/python3 %h/jarvis-voice/gpu-services/speaker_verify_service.py --threshold 0.45
Restart=on-failure
RestartSec=30
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
